# system/main/templates/pm2.service.erb
[Unit]
Description=PM2 next gen process manager for Node.js
After=network.target remote-fs.target

[Service]
Type=forking
User=developer
WorkingDirectory=<%= scope.lookupvar('settings::modulepath').split('/')[0 .. -2].join('/') %>

ExecStart=/usr/local/lib/node_modules/pm2/bin/pm2 start <%= scope.lookupvar('settings::modulepath').split('/')[0 .. -2].join('/') %>/cluster.json
ExecReload=/usr/local/lib/node_modules/pm2/bin/pm2 reload all
ExecStop=/usr/local/lib/node_modules/pm2/bin/pm2 kill

[Install]
WantedBy=multi-user.target
